# production-supervisor.conf
[program:juriidilisedjuured]
command=/home/juriidilisedjuured/juriidilisedjuured-backend/production-uwsgi.sh
autostart=true
autorestart=true
user=juriidilisedjuured
stopasgroup=true